# .config/systemd/user/wallpaper.timer
[Unit]
Description=changes the wallpaper every hour

[Timer]
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
